# task1/add_log_2_usb.rules
SUBSYSTEM=="usb",ACTION=="add",KERNEL=="sd[b-z][0-9]",RUN+="/bin/device_added.sh %k"
